# mqtt-s7-connector/rootfs/etc/services.d/gui/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.warning "Konfigurationsoberfläche wurde beendet (${1:-0})"
